# scripts/com.ascent.rpe-reminder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ascent.rpe-reminder</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jarvisforoli/projects/ascent/venv/bin/python</string>
        <string>/Users/jarvisforoli/projects/ascent/scripts/rpe_reminder.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jarvisforoli/projects/ascent/logs/rpe-reminder.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvisforoli/projects/ascent/logs/rpe-reminder.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/jarvisforoli/projects/ascent</string>
</dict>
</plist>
